# config/systemd/user/vm@.service
[Unit]
Description=Virtual Machine for: %i

[Service]
ExecStart=/usr/bin/VBoxHeadless -s %i
ExecStop=/usr/bin/VBoxManage controlvm %i savestate
